# roles/haproxy-for-workers/templates/scripts/etc/cron.d/ntaxa-haproxy-generate-config-from-workers

*/5 * * * * root { /usr/local/bin/ntaxa-artek-vz-update-haproxy.sh 2>&3 | logger -p user.notice; } 3>&1 1>&2 | logger -p user.error
